# image/files/etc/sv/ssh/run
#!/usr/bin/env svrun
# frozen_string_literal: true
# vim: ai ts=2 sts=2 et sw=2 ft=ruby

(ENV['SSH_PORT'] || 22).tap do |port|
  runner(['/usr/sbin/dropbear',
          '-Fsp', port,
          '-r', '/etc/dropbear/host_rsa']).call
end

# Local Variables:
# mode: ruby
# End:
