# image/files/etc/services/syslog/run
#!/usr/bin/env svrun
# vim: ai ts=2 sts=2 et sw=2 ft=ruby

# Send output to the container logs collector,
# send output to PID 1's stdio streams within the container.

service(%w[/sbin/syslogd -SnO /proc/1/fd/2]).call

# Local Variables:
# mode: ruby
# End:
